# xdg/autostart/set_wallpaper.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/set_default_kde_wallpaper-mx.sh
Hidden=false
NoDisplay=false
Name=Set Wallpaper
OnlyShowIn=KDE;
